# containers/etc/profile.d/aliases.sh
  echo "alias l='ls -lha'" >> ~/.bashrc && \
  echo "alias ll='ls -lha'" >> ~/.bashrc && \
  echo "alias la='ls -A'" >> ~/.bashrc && \
  echo "alias cls='clear && printf '\''\e[3J'\'''" >> ~/.bashrc 